# services/covidtracker.conf
[program:covid19tracker]
command=gunicorn --reload --workers=2 --bind 0.0.0.0:5000 app:app
directory=/usr/src/www
stdout_logfile=/tmp/logs/covid-19-tracker/covid.log
autostart=true
autorestart=true
redirect_stderr=True